# dist/myauthd/DEBIAN/postinst
systemctl enable myauthd
service myauthd start
